# scripts/etc/systemd/system/suspend-prevent.service
[Unit]
Description=Prevent suspend on lid close when on AC

[Service]
Type=simple
ExecStart=/opt/scripts/suspend-prevent.sh --forever
