# scripts/etc/systemd/system/suspend-prevent-atboot.service
[Unit]
Description=Determine during boot if inhibit lock needs to be set or not
After=dbus.service

[Service]
ExecStart=/opt/scripts/suspend-prevent.sh

[Install]
WantedBy=multi-user.target
