# udev/etc/udev/rules.d/50-suspend-prevent.rules
SUBSYSTEM=="power_supply", ENV{POWER_SUPPLY_ONLINE}=="0", RUN+="/opt/scripts/suspend-prevent.sh battery"
SUBSYSTEM=="power_supply", ENV{POWER_SUPPLY_ONLINE}=="1", RUN+="/opt/scripts/suspend-prevent.sh ac"
