# udev/etc/udev/rules.d/85-usb-headset.rules
KERNEL=="card?", SUBSYSTEM=="sound", ATTR{id}=="Headset", ACTION=="add", RUN+="/opt/scripts/usb-headset --fork"
