# systemd/.config/systemd/user/dropbox.service
[Unit]
Description=Dropbox

[Service]
ExecStart=/opt/bin/dropbox

[Install]
WantedBy=graphical.target
